# luci-app-netwiz/root/etc/init.d/netwiz-monitor
#!/bin/sh /etc/rc.common
# Copyright (C) 2026 huchd0 <https://github.com/huchd0/luci-app-netwiz>
# Licensed under the GNU General Public License v3.0

START=99
USE_PROCD=1
LOG_FILE="/tmp/netwiz.log"

log() {
    echo "$(date '+%F %T') [Service] $1" >> "$LOG_FILE"
    # 限制日志文件大小
    if [ $(wc -l < "$LOG_FILE" 2>/dev/null || echo 0) -gt 600 ]; then
        tail -n 500 "$LOG_FILE" > "$LOG_FILE.tmp"
        mv "$LOG_FILE.tmp" "$LOG_FILE"
    fi
}

start_service() {
    log "Netwiz 后台守护进程已启动"
    
    # 获取向导配置状态
    local configured=$(uci -q get netwiz.global.configured)
    
    # 判断密码是否为空
    local pwd_hash=$(awk -F: '/^root:/ {print $2}' /etc/shadow)
    local is_empty_pwd=0
    if [ -z "$pwd_hash" ] || [ "$pwd_hash" = "!" ] || [ "$pwd_hash" = "*" ] || [ "$pwd_hash" = "x" ]; then
        is_empty_pwd=1
    fi

    # 判断是否为固件自带安装
    local is_in_rom=0
    if [ -f /rom/etc/init.d/netwiz-monitor ]; then
        is_in_rom=1
    fi

    if [ "$configured" != "1" ]; then
        if [ "$is_empty_pwd" -eq 1 ] || [ "$is_in_rom" -eq 1 ]; then
            log "检测到首次开机，执行网关劫持"
            
            # 多语言适配：自动检测固件内置的本地语言
            local current_lang=$(uci -q get luci.main.lang)
            if [ "$current_lang" = "auto" ] || [ "$current_lang" = "en" ] || [ -z "$current_lang" ]; then
                # 过滤 internal 头部
                local smart_lang=$(uci -q show luci.languages 2>/dev/null | grep "luci\.languages\." | grep -v "\.en=" | head -n 1 | awk -F'=' '{print $1}' | awk -F'.' '{print $3}')
                
                # uci 里没有，直接扫描系统底层的汉化文件 (将 base.zh-cn.lmo 转换为 zh_cn)
                if [ -z "$smart_lang" ]; then
                    smart_lang=$(ls /usr/lib/lua/luci/i18n/base.*.lmo 2>/dev/null | grep -v "\.en\." | head -n 1 | awk -F'.' '{print $2}' | tr '-' '_')
                fi
    
                # 如果成功获取到了语言包（如 zh_cn, zh_tw）
                if [ -n "$smart_lang" ]; then
                    log "🔎 探测到本地语言包 [$smart_lang]，已自动应用为默认语言"
                    uci set luci.main.lang="$smart_lang"
                    uci commit luci
                    
                    # 清除 LuCI 页面缓存，强制下一次访问时重新渲染语言
                    rm -rf /tmp/luci-indexcache /tmp/luci-modulecache/* 2>/dev/null
                fi
            fi
            
            [ -f /www/cgi-bin/luci.bak ] || cp /www/cgi-bin/luci /www/cgi-bin/luci.bak
            
            # 写入分流劫持脚本
            if [ "$is_empty_pwd" -eq 1 ]; then
                cat << 'EOF' > /www/cgi-bin/luci
#!/bin/sh

# 1. 语言检测与缓存清理
IS_LANG_INIT=$(uci -q get netwiz.global.lang_init)
if [ "$IS_LANG_INIT" != "1" ]; then
    current_lang=$(uci -q get luci.main.lang)
    if [ "$current_lang" = "auto" ] || [ "$current_lang" = "en" ] || [ -z "$current_lang" ]; then
        smart_lang=$(ls /usr/lib/lua/luci/i18n/base.*.lmo 2>/dev/null | grep -v "\.en\." | head -n 1 | awk -F'.' '{print $2}' | tr '-' '_')
        if [ -n "$smart_lang" ]; then
            uci set luci.main.lang="$smart_lang"
            uci commit luci
            rm -rf /tmp/luci-indexcache /tmp/luci-modulecache/* 2>/dev/null
        fi
    fi
    uci -q set netwiz.global.lang_init='1'
    uci commit netwiz
fi

# 2. 页面拦截与向导跳转免密登录
CONF=$(uci -q get netwiz.global.configured)

if [ "$CONF" != "1" ]; then
    case "$PATH_INFO" in
        *rescue*)
            uci -q set netwiz.global.configured='1'
            uci commit netwiz
            [ -f /www/cgi-bin/luci.bak ] && mv -f /www/cgi-bin/luci.bak /www/cgi-bin/luci
            printf "Status: 302 Found\r\n"
            printf "Location: /cgi-bin/luci/\r\n"
            printf "Cache-Control: no-store, no-cache, must-revalidate, max-age=0\r\n\r\n"
            exit 0
            ;;
        *admin/netwiz*|*luci-static*|*rpc*)
            ;;
        *)
            # 输出 HTTP 头与隐形登录表单 (完美保留)
            printf "Status: 200 OK\r\n"
            printf "Content-Type: text/html; charset=utf-8\r\n"
            printf "Cache-Control: no-store, no-cache, must-revalidate, max-age=0\r\n\r\n"
            printf "<!DOCTYPE html><html><body style='background:#f8fafc;'>"
            printf "<form id='al' action='/cgi-bin/luci/admin/netwiz' method='post'>"
            printf "<input type='hidden' name='luci_username' value='root'>"
            printf "<input type='hidden' name='luci_password' value=''>"
            printf "</form><script>document.getElementById('al').submit();</script></body></html>\n"
            exit 0
            ;;
    esac
fi

exec /www/cgi-bin/luci.bak
EOF
            else
                cat << 'EOF' > /www/cgi-bin/luci
#!/bin/sh

# 1. 语言检测与缓存清理
IS_LANG_INIT=$(uci -q get netwiz.global.lang_init)
if [ "$IS_LANG_INIT" != "1" ]; then
    current_lang=$(uci -q get luci.main.lang)
    if [ "$current_lang" = "auto" ] || [ "$current_lang" = "en" ] || [ -z "$current_lang" ]; then
        smart_lang=$(ls /usr/lib/lua/luci/i18n/base.*.lmo 2>/dev/null | grep -v "\.en\." | head -n 1 | awk -F'.' '{print $2}' | tr '-' '_')
        if [ -n "$smart_lang" ]; then
            uci set luci.main.lang="$smart_lang"
            uci commit luci
            rm -rf /tmp/luci-indexcache /tmp/luci-modulecache/* 2>/dev/null
        fi
    fi
    uci -q set netwiz.global.lang_init='1'
    uci commit netwiz
fi

# 2. 页面拦截与向导跳转
CONF=$(uci -q get netwiz.global.configured)

if [ "$CONF" != "1" ];
then
    case "$PATH_INFO" in
        *rescue*)
            uci -q set netwiz.global.configured='1'
            uci commit netwiz
            [ -f /www/cgi-bin/luci.bak ] && mv -f /www/cgi-bin/luci.bak /www/cgi-bin/luci
            printf "Status: 302 Found\r\n"
            printf "Location: /cgi-bin/luci/\r\n"
            printf "Cache-Control: no-store, no-cache, must-revalidate, max-age=0\r\n\r\n"
            exit 0
            ;;
        *admin/netwiz*|*luci-static*|*rpc*)
            ;;
        *)
            # 已有密码则跳转手动登录
            printf "Status: 302 Found\r\n"
            printf "Location: /cgi-bin/luci/admin/netwiz\r\n"
            printf "Cache-Control: no-store, no-cache, must-revalidate, max-age=0\r\n\r\n"
            exit 0
            ;;
    esac
fi

exec /www/cgi-bin/luci.bak
EOF
            fi
            
            # 赋予执行权限
            chmod +x /www/cgi-bin/luci
        else
            log "检测到已有密码且为后期安装，取消向导"
            uci -q set netwiz.global.configured='1'
            uci commit netwiz
            
            # 恢复官方网关
            if [ -f /www/cgi-bin/luci.bak ]; then
                mv /www/cgi-bin/luci.bak /www/cgi-bin/luci
            fi
        fi
    else
        # 已配置过，确保恢复官方网关
        if [ -f /www/cgi-bin/luci.bak ]; then
            log "向导已完成，还原官方网关"
            mv /www/cgi-bin/luci.bak /www/cgi-bin/luci
        fi
    fi

    procd_open_instance
    procd_set_param command /usr/libexec/netwiz-monitor-loop.sh
    procd_set_param respawn
    procd_close_instance
}

stop_service() {
    log "Netwiz 后台守护进程已停止"
}
