# luci-app-netwiz/root/etc/init.d/netwiz-monitor
#!/bin/sh /etc/rc.common
# Copyright (C) 2026 huchd0 <https://github.com/huchd0/luci-app-netwiz>
# Licensed under the GNU General Public License v3.0
START=99
USE_PROCD=1

start_service() {
    procd_open_instance
    procd_set_param command /usr/libexec/netwiz-monitor-loop.sh
    procd_set_param respawn
    procd_close_instance
}
